# home/kiosk/.config/systemd/user/remote.service
[Unit]
Description=Console RFB Provider

[Service]
Type=simple
ExecStart=do-shadow-console
#Restart=always

[Install]
WantedBy=default.target
